# src/FDownl Docker/storage/fdownlstorage.service
[Unit]
Description=FDownl

[Service]
WorkingDirectory=/home/fdownl/app
ExecStart=/usr/bin/dotnet "/home/fdownl/app/Fdownl Storage.dll"
Restart=always

[Install]
WantedBy=multi-user.target
